# deploy/musicon.service
[Unit]
Description=Musicon API Server
After=network.target

[Service]
Type=simple
User=musicon
Group=musicon
WorkingDirectory=/opt/musicon
ExecStart=/opt/musicon/musicon-server
Restart=on-failure
RestartSec=5
Environment=SERVER_PORT=7847
Environment=DATABASE_URL=/opt/musicon/data/musicon.db
Environment=TJ_API_BASE_URL=https://www.tjmedia.com/legacy/api/newSongOfMonth

[Install]
WantedBy=multi-user.target
